# templates/etc/update-motd.d/10-custom-banner
#!/bin/sh
echo ""
echo "╔═══════════════════════════════════════════════════════════════════════════╗"
echo "║     ███████╗██╗   ██╗███████╗████████╗███████╗███╗   ███╗               ║"
echo "║     ██╔════╝╚██╗ ██╔╝██╔════╝╚══██╔══╝██╔════╝████╗ ████║               ║"
echo "║     ███████╗ ╚████╔╝ ███████╗   ██║   █████╗  ██╔████╔██║               ║"
echo "║     ╚════██║  ╚██╔╝  ╚════██║   ██║   ██╔══╝  ██║╚██╔╝██║               ║"
echo "║     ███████║   ██║   ███████║   ██║   ███████╗██║ ╚═╝ ██║               ║"
echo "║     ╚══════╝   ╚═╝   ╚══════╝   ╚═╝   ╚══════╝╚═╝     ╚═╝               ║"
echo "║                                                                           ║"
echo "║     ██████╗  █████╗ ████████╗██╗  ██╗         ™                         ║"
echo "║     ██╔══██╗██╔══██╗╚══██╔══╝██║  ██║                                   ║"
echo "║     ██████╔╝███████║   ██║   ███████║                                   ║"
echo "║     ██╔═══╝ ██╔══██║   ██║   ██╔══██║                                   ║"
echo "║     ██║     ██║  ██║   ██║   ██║  ██║                                   ║"
echo "║                                                                           ║"
echo "║                    Laravel Development Server                            ║"
echo "╚═══════════════════════════════════════════════════════════════════════════╝"
echo ""
echo "🎯 QUICK COMMANDS:"
echo "  server-help (shelp)  - Show all available commands"
echo "  server-info (sinfo)  - Show server status"
echo "  vdev                 - Start Vite dev server"
echo "  deploy               - Production build"
echo ""
echo "📁 Project: /var/www/html"
echo "===================================================="
echo ""